# build/spring/conf/spring.conf
[program:spring]
command=java -jar /opt/readus/readus-latest.jar
stderr_logfile = /var/log/supervisord/spring-stderr.log
stdout_logfile = /var/log/supervisord/spring-stdout.log
