# build/tornado/conf/tornado.conf
[program:tornado]
command=python3 /opt/readus/app.py
stderr_logfile = /var/log/supervisord/tornado-stderr.log
stdout_logfile = /var/log/supervisord/tornado-stdout.log
